# supervisord.conf
[supervisord]
nodaemon=true

[program:fdbserver]
command=/bin/bash /usr/lib/foundationdb/docker-start.sh
