# sql-server/supervisord.conf
[supervisord]
nodaemon=true

[program:fdbsql]
command=/bin/bash /usr/lib/foundationdb/docker-sql-layer.sh
